# deploy/systemd/schemaclient-gateway.service
# /etc/systemd/system/schemaclient-gateway.service
# Run Django ASGI app with websocket support via gunicorn + uvicorn worker.

[Unit]
Description=SchemaClient Tunnel Gateway (ASGI)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/schemaclient-gateway/schemaclienttunnel
EnvironmentFile=/home/ubuntu/schemaclient-gateway/.env
RuntimeDirectory=schemaclient-gateway
ExecStartPre=/usr/bin/mkdir -p /home/ubuntu/schemaclient-gateway/data
ExecStartPre=/usr/bin/bash -lc '/home/ubuntu/schemaclient-gateway/.venv/bin/python manage.py migrate --noinput'
ExecStartPre=/usr/bin/bash -lc '/home/ubuntu/schemaclient-gateway/.venv/bin/python manage.py collectstatic --noinput'
ExecStart=/home/ubuntu/schemaclient-gateway/.venv/bin/gunicorn gateway.asgi:application \
  --workers 1 \
  --worker-class uvicorn.workers.UvicornWorker \
  --bind unix:/home/ubuntu/schemaclient-gateway/gateway.sock \
  --umask 007 \
  --max-requests 2000 \
  --max-requests-jitter 200 \
  --timeout 120 \
  --graceful-timeout 30 \
  --keep-alive 75 \
  --access-logfile - \
  --error-logfile -
Restart=always
RestartSec=3
KillSignal=SIGQUIT
TimeoutStopSec=30
PrivateTmp=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
